# generators/app/templates/files/etc/profile.d/prompt.sh.erb
export PS1="\[$(tput setaf 1)\]\[$(tput bold)\][$NODE_ENV] \[$(tput sgr0)\]\u@\h:\W\\$ \[$(tput sgr0)\]"